# docker/supervisord.d/influxdb.conf
[program:influxdb]
command=/run.sh
stdout_logfile=/var/log/influxdb.log
redirect_stderr=true
