# docker/supervisord.d/statsd.conf
[program:statsd]
command=/usr/local/bin/statsd /config/statsd.conf
stdout_logfile=/var/log/statsd.log
redirect_stderr=true
